# live-build/config/includes.chroot/etc/systemd/system/sigmavault-api.service
[Unit]
Description=SigmaVault NAS OS API Server
Documentation=https://github.com/iamthegreatdestroyer/sigmavault-nas-os
After=network.target sigmavault-engined.service
Requires=sigmavault-engined.service

[Service]
Type=simple
User=sigmavault
Group=sigmavault
Environment=SIGMAVAULT_CONFIG=/etc/sigmavault/config.yaml
Environment=SIGMAVAULT_PORT=12080
Environment=SIGMAVAULT_RPC_URL=http://127.0.0.1:5000/api/v1
ExecStart=/opt/sigmavault/bin/sigmavault-api
Restart=always
RestartSec=5
StandardOutput=journal
StandardError=journal

# Security hardening
NoNewPrivileges=yes
ProtectSystem=strict
ProtectHome=yes
PrivateTmp=yes
ReadWritePaths=/var/lib/sigmavault /var/log/sigmavault

# Resource limits
MemoryMax=1G
CPUQuota=200%

[Install]
WantedBy=multi-user.target
